# modules/fireball/files/systemd/user/xorg-session.target
[Unit]
Description=User systemd services for xorg graphical sessio
BindsTo=graphical-session.target
